# root/etc/services.d/permissions/run
#!/bin/sh

sleep 20
chown -R unms:unms /home/app/unms/public
sleep 40