# root/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/cron -f -l 5