# root/etc/services.d/redis/run
#!/usr/bin/with-contenv sh

echo "Starting redis..."
s6-setuidgid abc redis-server --appendonly yes --dir /config/redis/
