# root/etc/services.d/rabbitmq/run
#!/usr/bin/with-contenv sh

export HOME=/var/lib/rabbitmq

echo "Starting rabbitmq-server..."
exec rabbitmq-server
