# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

echo "Starting nginx..."
s6-setuidgid root nginx -g 'daemon off;'
